# profile.d/linux-base-sgx.sh
#!/bin/sh
# Default to "out-of-proc" attestation
export SGX_AESM_ADDR=1
